# setup/unicorn/etc/systemd/system/unicorn.service
[Unit]
Description=Unicorn
After=network.target

[Service]
User=root
WorkingDirectory=/srv/http/benchmark
ExecStart=unicorn -E production --host [::] --port 80 --config-file /etc/unicorn.rb
Restart=on-failure

[Install]
WantedBy=multi-user.target
