# setup/passenger-standalone/etc/systemd/system/passenger.service
[Unit]
Description=Passenger Standalone
After=network.target

[Service]
User=root
WorkingDirectory=/srv/http/benchmark
ExecStart=passenger start --address '::' --port 80 --min-instances 8 --max-pool-size 8
Restart=on-failure

[Install]
WantedBy=multi-user.target
