# setup/cherrypy/etc/systemd/system/cherrypy.service
[Unit]
Description=CherryPy
After=network.target

[Service]
User=root
WorkingDirectory=/srv/http/benchmark
ExecStart=python server.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
