# setup/falcon/etc/systemd/system/falcon.service
[Unit]
Description=Falcon
After=network.target

[Service]
User=root
WorkingDirectory=/srv/http/benchmark
ExecStart=falcon serve --bind http://[::]:80
# This extends the concurrency model to a similar level as the default puma cluster:
# ExecStart=falcon serve --bind http://[::]:80 --hybrid --count 128
Restart=on-failure

[Install]
WantedBy=multi-user.target
