# setup/falcon/etc/systemd/system/falcon.service
[Unit]
Description=Falcon
After=network.target

[Service]
User=root
WorkingDirectory=/srv/http/benchmark
ExecStart=falcon serve --bind http://[::]:80
Restart=on-failure

[Install]
WantedBy=multi-user.target
