# setup/puma/etc/systemd/system/puma.service
[Unit]
Description=Puma
After=network.target

[Service]
User=root
WorkingDirectory=/srv/http/benchmark
ExecStart=puma --workers 8 --bind tcp://[::]:80
Restart=on-failure

[Install]
WantedBy=multi-user.target
